# jobs/clamav/templates/conf/clamav_logrotate.cron.erb
*/30 * * * * test -x /usr/sbin/logrotate && /usr/sbin/logrotate /var/vcap/jobs/clamav/conf/logrotate.conf >>/var/vcap/sys/log/clamav/clamav_logrotate_cron.log 2>&1