# jobs/clamav/templates/conf/clamav_dailyscan.cron.erb
30 01 * * * test -x /var/vcap/jobs/clamav/bin/clamav_scan.sh && /var/vcap/jobs/clamav/bin/clamav_scan.sh >>/var/vcap/sys/log/clamav/clamav.log 2>&1

